# raspbian-config/datas/etc/update-motd.d/10-uname
#!/bin/sh

_IP=$(/sbin/ifconfig eth0 | grep 'inet ' | awk '{ print $2}') || true
if [ "$_IP" ]; then
  printf "You can access to configuration here :\n http://%s\n\n\n" "$_IP"
fi
